# docker/php/supervisord.conf

[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm

[program:queue-worker]
command=php /var/www/html/artisan queue:work
